# content/etc/cron.d/pagerduty
* * * * * icinga /usr/local/bin/pagerduty_icinga.pl flush
